# airootfs/etc/systemd/system/reflector.service.d/archiso.conf
[Unit]
ConditionKernelCommandLine=!mirror

[Service]
Restart=on-failure
RestartSec=10
