# airootfs/etc/systemd/system/amd-perf-fix.service
[Unit]
Description=Set proper AMD GPU power profile

[Service]
ExecStart=/usr/bin/amd-perf-fix

[Install]
WantedBy=multi-user.target
